# src/data/bert-probe.service.in
[Unit]
Description=Telemetrics BERT Probe
Requires=telemprobd.socket
After=telemprobd.socket
ConditionPathExists=/etc/telemetrics/opt-in
ConditionPathExists=/sys/firmware/acpi/tables/data/BERT

[Service]
ExecStart=@bindir@/bertprobe
User=root

[Install]
WantedBy=multi-user.target
