# src/data/hprobe.service.in
[Unit]
Description=Telemetrics Heartbeat Probe
Requires=telemprobd.socket
After=telemprobd.socket
ConditionPathExists=/etc/telemetrics/opt-in

[Service]
ExecStart=@bindir@/hprobe -Hlu
User=telemetry

[Install]
WantedBy=multi-user.target
