# src/data/hprobe.service.in
[Unit]
Description=Telemetrics Heartbeat Probe
Requires=telemprobd.socket
After=telemprobd.socket
ConditionPathExists=!/etc/telemetrics/opt-out

[Service]
ExecStart=@bindir@/hprobe -H
User=telemetry

[Install]
WantedBy=multi-user.target
