# src/data/journal-probe-tail.service.in
[Unit]
Description=Telemetrics Systemd Journal Probe
Requires=telemprobd.socket
After=telemprobd.socket
ConditionPathExists=/etc/telemetrics/opt-in

[Service]
ExecStart=@bindir@/journalprobe -t
User=telemetry

[Install]
WantedBy=multi-user.target
