# src/data/journal-probe.service.in
[Unit]
Description=Telemetrics Systemd Journal Probe
Requires=telemd.socket
After=telemd.socket
ConditionPathExists=!/etc/telemetrics/opt-out

[Service]
ExecStart=@bindir@/journalprobe
User=telemetry

[Install]
WantedBy=multi-user.target
