# src/data/klogscanner.service.in
[Unit]
Description=Telemetrics Kernel Log Scanner
Requires=telemprobd.socket
After=telemprobd.socket oops-probe.service

[Service]
ExecStart=@bindir@/klogscanner

[Install]
WantedBy=multi-user.target
