# src/data/klogscanner.service.in
[Unit]
Description=Telemetrics Kernel Log Scanner
ConditionVirtualization=!container
Requires=telemprobd.socket
After=telemprobd.socket

[Service]
ExecStart=@bindir@/klogscanner

[Install]
WantedBy=multi-user.target
