# src/data/pstore-clean.service.in
[Unit]
Description=Telemetrics Pstore Clean Service

[Service]
ExecStart=@bindir@/pstoreclean

[Install]
WantedBy=multi-user.target
