# src/data/pstore-probe.service.in
[Unit]
Description=Telemetrics Pstore Probe
Requires=telemprobd.socket
After=telemprobd.socket
ConditionPathExists=/etc/telemetrics/opt-in

[Service]
ExecStart=@bindir@/pstoreprobe
User=telemetry

[Install]
WantedBy=multi-user.target
