# src/data/python-probe.service.in
[Unit]
Description=Telemetrics Python Probe
Requires=telemprobd.socket
After=telemprobd.socket
ConditionPathExists=!/etc/telemetrics/opt-out

[Service]
User=telemetry
ExecStart=@bindir@/pythonprobe

[Install]
WantedBy=multi-user.target

