# src/data/python-probe.path.in
[Unit]
Description=Python Probe Daemon Staging
ConditionPathExists=!/etc/telemetrics/opt-out

[Path]
DirectoryNotEmpty=/var/lib/telemetry/python

[Install]
WantedBy=multi-user.target
