# src/data/telemprobd-update-trigger.service.in
[Unit]
Description=Restart Telemetrics Daemon On Update
BindsTo=update-triggers.target

[Service]
Type=oneshot
ExecStart=@systemctldir@/systemctl try-restart telemprobd.service
