# src/data/telemprobd.service.in
[Unit]
Description=Telemetrics Daemon
ConditionPathExists=!/etc/telemetrics/opt-out

[Service]
ExecStart=@bindir@/telemprobd
WorkingDirectory=@localstatedir@/spool/telemetry
User=telemetry

[Install]
Also=telemprobd.socket
WantedBy=multi-user.target
