# src/data/telempostd.path.in
[Unit]
Description=Telemetry Post Daemon staging
ConditionPathExists=!/etc/telemetrics/opt-out

[Path]
DirectoryNotEmpty=@localstatedir@/spool/telemetry

[Install]
WantedBy=multi-user.target
