# src/data/telemprobd.socket.in
[Unit]
Description=Telemetrics Daemon
ConditionPathExists=/etc/telemetrics/opt-in

[Socket]
ListenStream=@SOCKETDIR@/telem-0

[Install]
WantedBy=sockets.target
